# share/doc/1wMonitoring/1wMonitoring.cron
#
# field          allowed values
# -----          --------------
# minute         0-59
# hour           0-23
# day of month   1-31
# month          1-12 (or names, see below)
# day of week    0-7 (0 or 7 is Sun, or use names)
# 

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * *	root	1wSensorMon.pl    2> /dev/null
*/5 * * * *	root	1wCounterMon.pl 2> /dev/null
